# supervisor/zookeeper.conf
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
